# scripts/files/hulk
#!/sbin/runscript

command="/home/nao/naoqi/bin/tuhhNao"
description="HULKs RoboCup Soccer SPL Software"
name="${description} 4.1 Ultimate Edition"
pidfile="/tmp/tuhhNao.pid"
outfile="/home/nao/naoqi/tuhhNao.out"
errfile="/home/nao/naoqi/tuhhNao.err"
user="nao"

depend() {
  need naoqi
}

restart() {
  stop
  sleep 2
  start
}

backup() {
  local FILE=$1
  local TIMESTAMP=$2
  if [ -f "$FILE" ]; then
    mv "$FILE" "$FILE.$TIMESTAMP"
  fi
  su -c "echo 'Starting $name at $TIMESTAMP' >> $FILE" nao
}

start() {
  local TIMESTAMP=$(date "+%Y-%m-%d_%H-%M-%S")
  backup "$outfile" "$TIMESTAMP"
  backup "$errfile" "$TIMESTAMP"
  ebegin "Starting $name"
    start-stop-daemon \
      --start --background \
      --exec "$command" \
      --user "$user" \
      --stdout "$outfile" --stderr "$errfile"
  eend $?
}

