# scripts/files/hulk
#!/sbin/runscript

command="/home/nao/naoqi/bin/tuhhNao"
description="HULKs RoboCup Soccer SPL Software"
name="${description} 4.1 Ultimate Edition"
pidfile="/tmp/tuhhNao.pid"
# The Gentoo version on the nao does not support the command_user variable :-(
start_stop_daemon_args="--user nao --background --stdout /tmp/tuhhNao.out --stderr /tmp/tuhhNao.err"

depend() {
  need naoqi
}

restart() {
  stop
  sleep 2
  start
}
